# systemd/openqa-consumer.service
[Unit]
Description=Run openQA tests in response to compose fedmsgs

[Service]
ExecStart=/usr/bin/fedora-openqa-consumer
Type=simple
